# WizardSpider/Conti Internal Software Leak/storage_ex-master-e4827b099abefd719fc674519ea0d2622ea304e0/need/storage_1_ex_2/logpost/supervisor/l_1300_1400.conf
[group:l_1300_1400]
programs=l_1300_1320,l_1320_1340,l_1340_1360,l_1360_1380,l_1380_1400

[program:l_1300_1320]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_1300_1320
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_1320_1340]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_1320_1340
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_1340_1360]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_1340_1360
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_1360_1380]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_1360_1380
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_1380_1400]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_1380_1400
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0