# WizardSpider/Conti Internal Software Leak/storage_ex-master-e4827b099abefd719fc674519ea0d2622ea304e0/need/storage_1_ex_2/logpost/supervisor/l_1800_1900.conf
[group:l_1800_1900]
programs=l_1800_1820,l_1820_1840,l_1840_1860,l_1860_1880,l_1880_1900

[program:l_1800_1820]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_1800_1820
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_1820_1840]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_1820_1840
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_1840_1860]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_1840_1860
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_1860_1880]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_1860_1880
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_1880_1900]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_1880_1900
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0