# WizardSpider/Conti Internal Software Leak/storage_ex-master-e4827b099abefd719fc674519ea0d2622ea304e0/need/storage_1_ex_2/logpost/supervisor/l_800_900.conf
[group:l_800_900]
programs=l_800_820,l_820_840,l_840_860,l_860_880,l_880_900

[program:l_800_820]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_800_820
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_820_840]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_820_840
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_840_860]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_840_860
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_860_880]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_860_880
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_880_900]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_880_900
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0