# WizardSpider/Conti Internal Software Leak/storage_ex-master-e4827b099abefd719fc674519ea0d2622ea304e0/need/storage_1_ex_2/logpost/supervisor/l_900_1000.conf
[group:l_900_1000]
programs=l_900_920,l_920_940,l_940_960,l_960_980,l_980_1000

[program:l_900_920]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_900_920
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_920_940]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_920_940
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_940_960]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_940_960
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_960_980]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_960_980
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_980_1000]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=logpost --name=l_980_1000
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0