# WizardSpider/Conti Internal Software Leak/srw-master-df4b6eddf7fdd2e07fb75d0492deeeb2e15f959e/need/supervisor/worker.conf
[program:worker]
directory=/root/srw
command=/usr/bin/php artisan worker
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0
