# WizardSpider/Conti Internal Software Leak/storage_ex-master-e4827b099abefd719fc674519ea0d2622ea304e0/need/storage_2/access/supervisor/a_0_100.conf
[group:a_0_100]
programs=a_0_20,a_20_40,a_40_60,a_60_80,a_80_100

[program:a_0_20]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:access --name=a_0_20
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:a_20_40]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:access --name=a_20_40
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:a_40_60]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:access --name=a_40_60
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:a_60_80]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:access --name=a_60_80
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:a_80_100]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:access --name=a_80_100
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0