# WizardSpider/Conti Internal Software Leak/storage_ex-master-e4827b099abefd719fc674519ea0d2622ea304e0/need/storage_1_ex/history/supervisor/h_0_100.conf
[group:h_0_100]
programs=h_0_20,h_20_40,h_40_60,h_60_80,h_80_100

[program:h_0_20]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:history --name=h_0_20
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:h_20_40]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:history --name=h_20_40
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:h_40_60]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:history --name=h_40_60
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:h_60_80]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:history --name=h_60_80
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:h_80_100]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:history --name=h_80_100
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0