# WizardSpider/Conti Internal Software Leak/storage_ex-master-e4827b099abefd719fc674519ea0d2622ea304e0/need/storage_1_ex/history/supervisor/h_200_300.conf
[group:h_200_300]
programs=h_200_220,h_220_240,h_240_260,h_260_280,h_280_300

[program:h_200_220]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:history --name=h_200_220
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:h_220_240]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:history --name=h_220_240
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:h_240_260]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:history --name=h_240_260
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:h_260_280]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:history --name=h_260_280
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:h_280_300]
directory=/root/storage_ex/
command=/usr/bin/php index.php --task=indexer:history --name=h_280_300
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0