# WizardSpider/Conti Internal Software Leak/storage_ex-master-e4827b099abefd719fc674519ea0d2622ea304e0/need/storage_1_ex/history/supervisor/searchd.conf
[group:history]
programs=history_main,history_0_100,history_100_200,history_200_300

[program:history_main]
command=/usr/bin/searchd --config /etc/sphinxsearch/storage_ex/history/searchd.conf --nodetach --iostats --cpustats
autostart=true
autorestart=true
numprocs=1
user=root

[program:history_0_100]
command=/usr/bin/searchd --config /etc/sphinxsearch/storage_ex/history/0_100.conf --nodetach --iostats --cpustats
autostart=true
autorestart=true
numprocs=1
user=root

[program:history_100_200]
command=/usr/bin/searchd --config /etc/sphinxsearch/storage_ex/history/100_200.conf --nodetach --iostats --cpustats
autostart=true
autorestart=true
numprocs=1
user=root

[program:history_200_300]
command=/usr/bin/searchd --config /etc/sphinxsearch/storage_ex/history/200_300.conf --nodetach --iostats --cpustats
autostart=true
autorestart=true
numprocs=1
user=root
