# WizardSpider/Conti Internal Software Leak/storage-master-3607d1f6a72e28efe84b55e8a660ff97db0e79a2/need/supervisor/servertstats_update_localsphinxstatus.conf
[program:serverstats_update_localsphinxstatus]
directory=/var/www/html
command=/usr/bin/php index.php --task=serverstats:update:localSphinxStatus
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
