# WizardSpider/Conti Internal Software Leak/storage-master-3607d1f6a72e28efe84b55e8a660ff97db0e79a2/need/supervisor/storage_1.conf
[program:searchd]
command=/usr/bin/searchd --config /etc/sphinxsearch/storage/storage_1.conf --nodetach --iostats --cpustats
autostart=true
autorestart=true
numprocs=1
user=root
