# WizardSpider/Conti Internal Software Leak/storage_ex-master-e4827b099abefd719fc674519ea0d2622ea304e0/need/storage_cutted_ex/logpost/supervisor/l_300_400.conf
[group:l_300_400]
programs=l_300_320,l_320_340,l_340_360,l_360_380,l_380_400

[program:l_300_320]
directory=/root/storage_ex/
command=/usr/bin/php5.6 index.php --task=logpost --name=l_300_320
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_320_340]
directory=/root/storage_ex/
command=/usr/bin/php5.6 index.php --task=logpost --name=l_320_340
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_340_360]
directory=/root/storage_ex/
command=/usr/bin/php5.6 index.php --task=logpost --name=l_340_360
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_360_380]
directory=/root/storage_ex/
command=/usr/bin/php5.6 index.php --task=logpost --name=l_360_380
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0

[program:l_380_400]
directory=/root/storage_ex/
command=/usr/bin/php5.6 index.php --task=logpost --name=l_380_400
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0